# base_portal/roles/portal/templates/scripts/crontab/andino-crontab.j2
0 0 * * * /etc/ckan_init.d/clean_datastore.sh
0 0 * * * /etc/ckan_init.d/update_datastore.sh
